# etc/init.d/battery_watchdog
#!/sbin/openrc-run

description="Checks battery status/level every 30 secs. and if it is at a critical state, displays a splash and turns off the device."

depend() {
	keyword -prefix -lxc
}

battery_check() {
	DEVICE="$(cat /opt/inkbox_device)"
	DEVICE_VARIANT="$(uname -r | grep -o '[^-]*$')"

	if [ "${DEVICE}" == "kt" ]; then
		SYSFS_PATH_CAPACITY="/sys/devices/system/yoshi_battery/yoshi_battery0/battery_capacity"
		SYSFS_PATH_STATUS="/sys/devices/system/yoshi_battery/yoshi_battery0/battery_status"
		STATUS_CHARGING=1
	elif [ "${DEVICE}" == "n306" ] && [ "${DEVICE_VARIANT}" == "n306c" ]; then
		SYSFS_PATH_CAPACITY="/sys/class/power_supply/battery/capacity"
		SYSFS_PATH_STATUS="/sys/class/power_supply/battery/status"
		STATUS_CHARGING="Charging"
	elif [ "${DEVICE}" == "n249" ]; then
		SYSFS_PATH_CAPACITY="/sys/class/power_supply/rn5t618-battery/capacity"
		SYSFS_PATH_STATUS="/sys/class/power_supply/rn5t618-battery/status"
		STATUS_CHARGING="Charging"
	else
		SYSFS_PATH_CAPACITY="/sys/devices/platform/pmic_battery.1/power_supply/mc13892_bat/capacity"
		SYSFS_PATH_STATUS="/sys/devices/platform/pmic_battery.1/power_supply/mc13892_bat/status"
		STATUS_CHARGING="Charging"
	fi

	FBINK="/opt/bin/fbink/fbink"
	IMAGES_PATH="/etc/images"

	if grep -q "${STATUS_CHARGING}" "${SYSFS_PATH_STATUS}" &>/dev/null; then
		# Skip since the battery is charging
		:
	else
		BT_LEVEL=$(cat "${SYSFS_PATH_CAPACITY}")
		if [ ${BT_LEVEL} -le 5 ]; then
			"${FBINK}" -k -f -q
			sleep 0.5
			"${FBINK}" -k -f -q
			sleep 0.5
			"${FBINK}" -q -g file="${IMAGES_PATH}/${DEVICE}-battery_critical.png"
			sleep 0.5
			"${FBINK}" -q -s -f
			# Prevent other processes from printing over what we just showed on the screen
			rm -f /dev/fb0
			/sbin/poweroff no_splash
		else
			:
		fi
	fi
}

start() {
	ebegin "Starting battery watchdog"
	DEVICE="$(cat /opt/inkbox_device)"
	DEVICE_VARIANT="$(uname -r | grep -o '[^-]*$')"

	if [ "${DEVICE}" != "emu" ]; then
		while true; do
			if grep -q "true" "/kobo/tmp/inkbox_running" &>/dev/null; then
				if grep -q "true" "/tmp/sleep_mode" &>/dev/null; then
					battery_check
				else
					:
				fi
			else
				battery_check
			fi
			sleep 30
		done &
	else
		einfo "Not starting battery watchdog since InkBox is running in an emulator."
	fi
	eend $?
}
