# etc/init.d/encfs
#!/sbin/openrc-run

description="Sets up encrypted storage."

depend() {
	keyword -prefix -lxc
}

error_exit() {
	echo "Error: ${1}. Exiting."
	[ -z "${2}" ] && exit 1
	exit "${2}"
}

calculate() {
	awk "BEGIN { print "$*" }";
}

start() {
	BOOTSTRAP_CREATE=$(cat /run/encfs/encrypted_storage_create 2>/dev/null)
	if ! echo "${BOOTSTRAP_CREATE}" | grep -q "true"; then
		ENCRYPTED_STORAGE_CONFIG=$(cat /opt/config/18-encrypted_storage/status 2>/dev/null)
		STORAGE_LIST=$(cat /opt/config/18-encrypted_storage/storage_list 2>/dev/null)
		if ! echo "${ENCRYPTED_STORAGE_CONFIG}" | grep -q "true"; then echo "Storage encryption disabled. Exiting."; exit 0; fi
		ENCRYPTED_STORAGE_ARCHIVE=$(cat /run/encfs/encrypted_storage_archive 2>/dev/null)
		ENCRYPTED_STORAGE_MOUNTPOINT=$(cat /run/encfs/encrypted_storage_mountpoint 2>/dev/null)
		ENCRYPTED_STORAGE_PASSPHRASE=$(cat /run/encfs/encrypted_storage_passphrase 2>/dev/null)
		if [ ! -z "${ENCRYPTED_STORAGE_ARCHIVE}" ] && [ ! -z "${ENCRYPTED_STORAGE_MOUNTPOINT}" ] && [ ! -z "${ENCRYPTED_STORAGE_PASSPHRASE}" ]; then
			# Encrypted storage archive is a SquashFS file, mounting it
			mkdir -p /run/encfsd/encrypted
			squashfuse "${ENCRYPTED_STORAGE_ARCHIVE}" /run/encfsd/encrypted
			if ! echo "${STORAGE_LIST}" | grep -q "${ENCRYPTED_STORAGE_MOUNTPOINT}"; then
				error_exit "Mountpoint is not in trusted storage list" 2
			else
				# Writing script which returns passphrase to encfs program
				printf "#!/bin/sh\nprintf '${ENCRYPTED_STORAGE_PASSPHRASE}\p'\nexit 0\n" | sed 's/\\p/\\n/g' > /run/encfs/encfs_passphrase_return.sh && chmod +x /run/encfs/encfs_passphrase_return.sh
				encfs --extpass=/run/encfs/encfs_passphrase_return.sh /run/encfsd/encrypted "${ENCRYPTED_STORAGE_MOUNTPOINT}" || error_exit "Encrypted storage setup" 1
				echo "true" > /run/encrypted_storage_setup
			fi
			exit 0
		else
			error_exit "Arguments parsing" 2
		fi
	else
		BOOTSTRAP_FILES_LOCATION=$(cat /run/encfs/encrypted_storage_bootstrap_files_location 2>/dev/null)
		BOOTSTRAP_ARCHIVE_LOCATION=$(cat /run/encfs/encrypted_storage_bootstrap_archive_location 2>/dev/null)
		BOOTSTRAP_PASSPHRASE=$(cat /run/encfs/encrypted_storage_bootstrap_passphrase 2>/dev/null)
		MEMTOTAL=$(grep "MemTotal" /proc/meminfo | awk '{ print $2 }')
		[ ! -e "${BOOTSTRAP_FILES_LOCATION}" ] && error_exit "Bootstrap files location not found" 127
		[ -z "${BOOTSTRAP_ARCHIVE_LOCATION}" ] && BOOTSTRAP_ARCHIVE_LOCATION="${BOOTSTRAP_FILES_LOCATION}/data.encfs"
		BOOTSTRAP_FILES_LOCATION_FOLDER_SIZE=$(du -s "${BOOTSTRAP_FILES_LOCATION}" | awk '{ print $1 }')
		if [ -z "${BOOTSTRAP_FILES_LOCATION}" ] || [ -z "${BOOTSTRAP_PASSPHRASE}" ]; then
			error_exit "Arguments parsing" 127
		fi
		AVAILABLE_RAM=$(calculate ${MEMTOTAL}/4)
		BOOTSTRAP_FILES_LOCATION_FOLDER_SIZE=$(calculate ${BOOTSTRAP_FILES_LOCATION_FOLDER_SIZE}*2)

		if [ ${BOOTSTRAP_FILES_LOCATION_FOLDER_SIZE} -gt ${AVAILABLE_RAM} ]; then
			BOOTSTRAP_ROOT_ENCRYPTED="/data/storage/encfsd/bootstrap_root_encrypted"
			BOOTSTRAP_ROOT_VISIBLE="/data/storage/encfsd/bootstrap_root_visible"
		else
			BOOTSTRAP_ROOT_ENCRYPTED="/run/encfsd/bootstrap_root_encrypted"
			BOOTSTRAP_ROOT_VISIBLE="/run/encfsd/bootstrap_root_visible"
		fi
		mkdir -p "${BOOTSTRAP_ROOT_ENCRYPTED}"
		mkdir -p "${BOOTSTRAP_ROOT_VISIBLE}"

		# Writing script which returns passphrase to encfs program
		printf "#!/bin/sh\nprintf '${BOOTSTRAP_PASSPHRASE}\p'\nexit 0\n" | sed 's/\\p/\\n/g' > /run/encfs/encfs_bootstrap_return.sh && chmod +x /run/encfs/encfs_bootstrap_return.sh
		encfs --paranoia --extpass=/run/encfs/encfs_bootstrap_return.sh "${BOOTSTRAP_ROOT_ENCRYPTED}" "${BOOTSTRAP_ROOT_VISIBLE}" || error_exit "Encrypted storage setup" 2
		rm -rf /run/encfs
		find "${BOOTSTRAP_FILES_LOCATION}" -type f -exec cp -r {} "${BOOTSTRAP_ROOT_VISIBLE}" \; || error_exit "Files copy" 1
		mksquashfs "${BOOTSTRAP_ROOT_ENCRYPTED}" "${BOOTSTRAP_ARCHIVE_LOCATION}" -b 1048576 -comp xz -Xdict-size 100% -always-use-fragments || error_exit "Failure in compressing encrypted archive" 1
		sync
		umount -l -f "${BOOTSTRAP_ROOT_VISIBLE}"
		rm -rf /run/encfsd /data/storage/encfsd
		sync
		echo "true" > /run/encrypted_storage_bootstrap_setup
	fi
}
