# etc/init.d/evtest_watch
#!/sbin/openrc-run

description="Redirects evtest input events list to files in a tmpfs"

depend() {
	before sleep_standby
	before inkbox_gui
	keyword -prefix -lxc
}

start() {
	ebegin "Starting evtest watch"
	evtest --grab /dev/input/event0 > /tmp/power &
	evtest /dev/input/event1 > /tmp/touch_event &
	eend ${?}
}

stop() {
	ebegin "Stopping evtest watch"
	killall -q evtest
	eend ${?}
}
