# etc/init.d/ipd
#!/sbin/openrc-run

description="Launches InkBox Power Daemon (ipd)"
DEVICE="$(cat /opt/inkbox_device)"

depend() {
	keyword -prefix -lxc
	before inkbox_gui gui_apps
	after modules_init usbnet splash onboard_mount
}

turn_off_led() {
	if [ "${DEVICE}" == "n705" ] || [ "${DEVICE}" == "n905b" ] || [ "${DEVICE}" == "n905c" ] || [ "${DEVICE}" == "n613" ]; then
		echo 0 > "/sys/class/leds/pmic_ledsg/brightness"
	elif [ "${DEVICE}" == "n306" ]; then
		echo 1 > "/sys/devices/platform/leds/leds/GLED/brightness" ; echo 0 > "/sys/devices/platform/leds/leds/GLED/brightness"
	elif [ "${DEVICE}" == "n873" ] || [ "${DEVICE}" == "n418" ]; then
		echo 1 > "/sys/class/leds/GLED/brightness" ; echo 0 > "/sys/class/leds/GLED/brightness"
	elif [ "${DEVICE}" == "bpi" ]; then
		echo 0 > "/sys/devices/platform/leds/leds/bpi:red:pwr/brightness"
	elif [ "${DEVICE}" == "n236" ] || [ "${DEVICE}" == "n437" ]; then
		/opt/bin/shutdown_led
	elif [ "${DEVICE}" == "n249" ]; then
		echo 0 > "/sys/class/leds/e60k02:white:on/brightness"
	elif [ "${DEVICE}" == "emu" ]; then
		:
	else
		echo 0 > "/sys/class/leds/pmic_ledsg/brightness"
	fi
}

start() {
	ebegin "Starting InkBox Power Daemon (ipd)"
	killall -q -9 evtest
	turn_off_led
	rename-input-devices &
	DEBUG="$(cat /boot/flags/IPD_DEBUG 2>/dev/null)" /sbin/ipd &
	eend ${?}
}

stop() {
	ebegin "Stopping InkBox Power Daemon (ipd)"
	killall -9 ipd
	eend ${?}
}
