# etc/init.d/mass_storage
#!/sbin/openrc-run

description="Uses the g_mass_storage module to add mass storage capabilities to the device"

depend() {
	keyword -prefix -lxc
}

start() {
	ebegin "Exporting onboard storage over USB..."
	FILE=$(cat /run/mass_storage_export)
	DEVICE=$(cat /opt/inkbox_device)
	[ "$(uname -r | grep -o '[^-]*$')" == "n306c" ] && DEVICE_VARIANT="n306c"

	if [ -z "${FILE}" ]; then
		FILE="/opt/storage/onboard"
	fi
	rc-service usbnet stop
	if [ "${DEVICE}" == "n306" ] || [ "${DEVICE}" == "n249" ] || [ "${DEVICE}" == "n873" ]; then
		if [ "${DEVICE_VARIANT}" != "n306c" ]; then
			insmod /modules/fs/configfs/configfs.ko
			insmod /modules/drivers/usb/gadget/libcomposite.ko
			insmod /modules/drivers/usb/gadget/function/usb_f_mass_storage.ko
		fi
	fi
	if [ "${DEVICE}" != "emu" ]; then
		insmod /modules/g_mass_storage.ko file="${FILE}" removable=y stall=0
	fi
	eend $?
}
stop() {
	ebegin "Exiting USBMS session..."
	if [ "${DEVICE}" == "n306" ] || [ "${DEVICE}" == "n249" ] || [ "${DEVICE}" == "n873" ]; then
		rmmod /modules/drivers/usb/gadget/function/usb_f_mass_storage.ko
		rmmod /modules/drivers/usb/gadget/libcomposite.ko
		rmmod /modules/fs/configfs/configfs.ko
	fi
	if [ "${DEVICE}" != "emu" ]; then
		rmmod /modules/g_mass_storage.ko
	fi
	eend $?
}
