# etc/init.d/mass_storage
#!/sbin/openrc-run

description="Uses the g_mass_storage module to add mass storage capabilities to the device"

depend() {
	keyword -prefix -lxc
}

start() {
	ebegin "Exporting onboard storage over USB..."
	FILE=$(cat /run/mass_storage_export)
	DEVICE=$(cat /opt/inkbox_device)
	if [ -z "${FILE}" ]; then
		FILE="/opt/storage/onboard"
	fi
	rc-service usbnet stop
	if [ "${DEVICE}" != "emu" ]; then
		insmod /modules/g_mass_storage.ko file="${FILE}" removable=y stall=0
	fi
	eend $?
}
stop() {
	ebegin "Exiting USBMS session..."
	if [ "${DEVICE}" != "emu" ]; then
		rmmod /modules/g_mass_storage.ko
	fi
	eend $?
}
