# etc/init.d/modules_init
#!/sbin/openrc-run

description="Loads necessary core kernel modules."

depend() {
	before mass_storage usbnet
	keyword -prefix -lxc
}

start() {
	DEVICE=$(cat /opt/inkbox_device)
	if [ "${DEVICE}" != "n873" ] && [ "${DEVICE}" != "n236" ] && [ "${DEVICE}" != "n437" ] || [ "${DEVICE}" != "emu" ] && [ "${DEVICE}" != "bpi" ]; then
		insmod /modules/arcotg_udc.ko
	fi
	eend $?
}
