# etc/init.d/sleep_standby
#!/sbin/openrc-run

description="Waits for a power button event and SIGSTOPs InkBox"

depend() {
	before inkbox_gui
	keyword -prefix -lxc
}

start() {
	ebegin "Waiting for a power button event to put the device to sleep"
	/etc/init.d/sleep_standby.sh &
	eend $?
}
