# etc/init.d/usbnet
#!/sbin/openrc-run

description="Starts USBNetworking"

depend() {
	keyword -prefix -lxc
}

start() {
	ebegin "Starting USBNet"
	ROOTED=$(cat /opt/root/rooted)
	FLAG=$(cat /boot/flags/USBNET_ENABLE 2>/dev/null)
	DEVICE=$(cat /opt/inkbox_device)

	if [ "${ROOTED}" == "true" ]; then
		if [ "${FLAG}" == "true" ]; then
			rc-service mass_storage stop &>/dev/null
			if [ "${DEVICE}" == "n705" ] || [ "${DEVICE}" == "n905b" ] || [ "${DEVICE}" == "n905c" ] || [ "${DEVICE}" == "n613" ] || [ "${DEVICE}" == "n236" ] || [ "${DEVICE}" == "n437" ]; then
				insmod "/modules/g_ether.ko"
			elif [ "${DEVICE}" == "n873" ] || [ "${DEVICE}" == "bpi" ]; then
				insmod "/modules/fs/configfs/configfs.ko"
				insmod "/modules/drivers/usb/gadget/libcomposite.ko"
				insmod "/modules/drivers/usb/gadget/function/u_ether.ko"
				insmod "/modules/drivers/usb/gadget/function/usb_f_ecm.ko"
				[ -e "/modules/drivers/usb/gadget/function/usb_f_eem.ko" ] && insmod "/modules/drivers/usb/gadget/function/usb_f_eem.ko"
				insmod "/modules/drivers/usb/gadget/function/usb_f_ecm_subset.ko"
				insmod "/modules/drivers/usb/gadget/function/usb_f_rndis.ko"
				insmod "/modules/drivers/usb/gadget/legacy/g_ether.ko"
			elif [ "${DEVICE}" == "emu" ]; then
				:
			else
				insmod "/modules/g_ether.ko"
			fi
			ifconfig usb0 up
			ifconfig usb0 192.168.2.2
		else
			ewarn "USB networking is disabled. Please set the USBNET_ENABLE flag and try again."
		fi
	else
		# Device is not rooted; enforcing security policy.
		eerror "Permission denied; USB networking is disabled."
	fi
	eend $?
}
stop() {
	ebegin "Stopping USBNet"
	if [ "${DEVICE}" != "emu" ]; then
		rmmod /modules/g_ether.ko
	fi
	eend $?
}
